# rules/99-dvd-monitor.rules
ACTION=="change", SUBSYSTEM=="block", KERNEL=="sr[0-9]*", RUN+="/usr/local/bin/dwmdvd $env{DEVNAME}"
